# deb/DEBIAN/postinst
#!/bin/bash

update-rc.d pwr-and-control-button-monitor defaults
update-rc.d led-daemon defaults
